# terraform/apps/nginx/scripts/token_mgmt.service
[Unit]
Description=Token Management script for Vault & Consul-template

[Service]
Type=simple
ExecStart=/usr/local/bin/token_mgmt.sh
